# foreman-kiosk.service
[Unit]
Description=Shop Schedule Kiosk
After=graphical.target network.target

[Service]
User=__USER__
Environment=DISPLAY=:0
Environment=XAUTHORITY=%h/.Xauthority
ExecStartPre=/bin/sleep 5
ExecStart=/usr/bin/chromium \
  --kiosk \
  --noerrdialogs \
  --disable-infobars \
  --no-first-run \
  --disable-session-crashed-bubble \
  --disable-restore-session-state \
  http://localhost:8080/kiosk.html
Restart=always
RestartSec=5

[Install]
WantedBy=graphical.target
